# roles/nextcloud/files/nextcloud_ldap_cleanup@.timer
[Unit]
Description=Nextcloud's LDAP cleanup [%i]

[Timer]
OnCalendar=*-*-* 00:10:00
Persistent=False

[Install]
WantedBy=timers.target
